# recipes-support/fastrpc/fastrpc/adsprpcd.service
[Unit]
Description=aDSP RPC daemon

[Service]
Type=exec
ExecStart=/usr/bin/adsprpcd

[Install]
WantedBy=multi-user.target
