# recipes-support/fastrpc/fastrpc/cdsprpcd.service
[Unit]
Description=cDSP RPC daemon

[Service]
Type=exec
ExecStart=/usr/bin/cdsprpcd

[Install]
WantedBy=multi-user.target
