# recipes-support/fastrpc/fastrpc/sdsprpcd.service
[Unit]
Description=sDSP RPC daemon

[Service]
Type=exec
ExecStart=/usr/bin/sdsprpcd

[Install]
WantedBy=multi-user.target
