# recipes-support/qbootctl/files/qbootclt-bless-boot.service.in
#  SPDX-License-Identifier: LGPL-2.1+
#
#  This file was copied from systemd at [1]
#
#  https://github.com/systemd/systemd/blob/main/units/systemd-bless-boot.service.in
#
#  systemd is free software; you can redistribute it and/or modify it
#  under the terms of the GNU Lesser General Public License as published by
#  the Free Software Foundation; either version 2.1 of the License, or
#  (at your option) any later version.

[Unit]
Description=Mark the Current Boot slot as Good
DefaultDependencies=no
Requires=boot-complete.target
After=local-fs.target boot-complete.target
Conflicts=shutdown.target
Before=shutdown.target

[Service]
Type=oneshot
RemainAfterExit=yes
ExecStart=@bindir@/qbootctl -m
